# roles/iperf/templates/iperf.service.j2
[Unit]
Description=iperf - network measurement tool
After=network.target

[Service]
ExecStart=/usr/bin/iperf -s -w {{ iperf_tcp_window_size }}

[Install]
WantedBy=multi-user.target